# scripts/prod/system_services/systemd/airflow-api.service
# =============================================================================
# Airflow API Server (Web UI)
# =============================================================================

[Unit]
Description=Airflow API Server (Web UI)
After=network.target

[Service]
Type=simple
User=etluser
Group=etluser
WorkingDirectory=/home/etluser/ai-projects/data_airflow_server
EnvironmentFile=/home/etluser/ai-projects/data_airflow_server/.env.prod

ExecStart=/bin/bash -c 'source /home/etluser/ai-projects/data_airflow_server/.venv/bin/activate && airflow api-server'

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target