# scripts/prod/system_services/systemd/airflow-scheduler.service
# =============================================================================
# scripts/prod/system_services/systemd/airflow-scheduler.service
# Airflow Scheduler - Clean modern version
# =============================================================================

[Unit]
Description=Airflow Scheduler
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=etluser
Group=etluser
WorkingDirectory=${PROJECT_ROOT}
Environment=PROJECT_ROOT=/home/etluser/ai-projects/data_airflow_server
EnvironmentFile=${PROJECT_ROOT}/.env.prod

ExecStartPre=/bin/bash -c 'source .venv/bin/activate && airflow dags reserialize'
ExecStart=/bin/bash -c 'source .venv/bin/activate && airflow scheduler'

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target